# deploy/systemd/asr-pool.service
[Unit]
Description=ASR Pool (Uvicorn)

[Service]
Type=simple
WorkingDirectory=%h/projects/asr-pool
ExecStart=%h/projects/asr-pool/.venv/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 18090
Restart=always
RestartSec=2
EnvironmentFile=-%h/.config/asr-pool/asr-pool.env

[Install]
WantedBy=default.target
