# deploy/systemd/transcribe-asr-pool-dev.service
[Unit]
Description=ASR Pool Dev (Uvicorn)

[Service]
Type=simple
# Replace <asr-pool-dir> with the absolute path to your local checkout.
WorkingDirectory=<asr-pool-dir>
ExecStart=<asr-pool-dir>/.venv/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 18090
Restart=always
RestartSec=2
EnvironmentFile=-%h/.config/asr-pool-dev/dev.env

[Install]
WantedBy=default.target
